# packaging/linux/debuntu/lenny/aghdvic/debian/postinst
#!/bin/sh

echo " POSTINST: Resetting network buffers with sysctl"
invoke-rc.d procps start

