# packaging/linux/slackware/SL-11.0/mDNSResponder/rc.local

# Startup the mDNSResponder daemon
#
if [ -x /etc/rc.d/rc.mdnsd ]; then
  /etc/rc.d/rc.mdnsd restart
fi

